# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:functions]
command=/usr/local/bin/deno run --allow-net --allow-env --allow-read /app/functions/server.ts
directory=/app/functions
autostart=true
autorestart=true
priority=20
environment=
    SUPABASE_URL="%(ENV_SUPABASE_URL)s",
    SUPABASE_ANON_KEY="%(ENV_SUPABASE_ANON_KEY)s",
    SUPABASE_SERVICE_ROLE_KEY="%(ENV_SUPABASE_SERVICE_ROLE_KEY)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    STRIPE_SECRET_KEY="%(ENV_STRIPE_SECRET_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
